# root/etc/services.d/apache/run
#!/usr/bin/with-contenv sh

/usr/local/bin/run.sh

